# files/system/etc/profile.d/hardened_malloc.sh
#!/usr/bin/sh
export LD_PRELOAD='libno_rlimit_as.so'
