# files/system/etc/profile.d/hardened_malloc.sh
#!/usr/bin/sh
export LD_PRELOAD='libmimalloc.so.2 libno_rlimit_as.so'
